# base/.config/autostart/picom.desktop
# Disable picom autostart, as it is already started by xinit
[Desktop Entry]
Type=Application
Name=picom
TryExec=picom
Exec=picom
Hidden=true
